# deploy/backup.timer
[Unit]
Description=Roda backup.service diariamente às 03:00 (America/Sao_Paulo)

[Timer]
OnCalendar=*-*-* 03:00:00 America/Sao_Paulo
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target
